# scripts/com.algotrading.king.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.algotrading.king.plist
  Auto-start & auto-restart the trading bot (paper mode).

  INSTALL:
    cp scripts/com.algotrading.king.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.algotrading.king.plist

  UNINSTALL:
    launchctl unload ~/Library/LaunchAgents/com.algotrading.king.plist
    rm ~/Library/LaunchAgents/com.algotrading.king.plist

  NOTES:
  - Always starts in PAPER mode (safe for auto-restart).
  - For live trading, use the controlled launcher: python3 scripts/go_live.py
  - Logs stdout to logs/service/bot.log, stderr to logs/service/bot_error.log
  - KeepAlive = true → restarts automatically if the process exits.
  - RunAtLoad = true → starts when you load the plist (and on reboot).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.algotrading.king</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>-B</string>
        <string>/Users/joshmacbookair2020/Projects/algo_trading_final/scripts/boot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/joshmacbookair2020/Projects/algo_trading_final</string>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start on load (and on reboot via LoginItems) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Wait 30 seconds before restarting after a crash (prevents .pyc lock contention) -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Capture logs -->
    <key>StandardOutPath</key>
    <string>/Users/joshmacbookair2020/Projects/algo_trading_final/logs/service/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joshmacbookair2020/Projects/algo_trading_final/logs/service/bot_error.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/joshmacbookair2020</string>
        <!-- Prevent Python 3.14 .pyc file lock deadlock on fast restarts -->
        <key>PYTHONDONTWRITEBYTECODE</key>
        <string>1</string>
        <!-- Disable tqdm progress bars — prevents OSError [Errno 11] deadlock in launchd (no TTY) -->
        <key>TQDM_DISABLE</key>
        <string>1</string>
        <key>TOKENIZERS_PARALLELISM</key>
        <string>false</string>
    </dict>
</dict>
</plist>
